# infra/cloudflared-api.service
[Unit]
Description=Cloudflare Quick Tunnel — Kubera API
After=network.target
Requires=docker.service
After=docker.service

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/bin/cloudflared tunnel --url http://localhost:3002
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
